# docker/ldap-server/etc_supervisord.d/hrzoo-ldap.ini
[supervisord]
nodaemon=true

[program:openldap]
command=/usr/sbin/slapd -u ldap -h "ldap:/// ldaps:/// ldapi:///" -d 1
redirect_stderr=true
stdout_logfile=/home/user/ldap-logs
autorestart=true
